# system/freezerbot-updater.timer
[Unit]
Description=Run Freezerbot firmware updater nightly at 3am

[Timer]
# Run every night at 3am
OnCalendar=*-*-* 03:00:00
# Randomize start time within a 10-minute window to avoid exact timing collisions
RandomizedDelaySec=600s
# Set how accurately systemd should try to hit the scheduled time
AccuracySec=1m
# Remember last trigger time across reboots
Persistent=true

[Install]
WantedBy=timers.target